# launchagent/ai.mirrordna.brain.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.mirrordna.brain.api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mirror-admin/repos/mirrorbrain-api/.venv/bin/uvicorn</string>
        <string>src.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8100</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mirror-admin/repos/mirrorbrain-api</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mirror-admin/.mirrordna/brain/api.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mirror-admin/.mirrordna/brain/api.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MIRRORBRAIN_DATA_PATH</key>
        <string>/Users/mirror-admin/.mirrordna/brain/data</string>
    </dict>
</dict>
</plist>
